# etc/crontab
* *     * * *   root    /home/balkonkraftwerk/balkonkraftwerk/check_tibber.py || systemctl restart tibber_pulse.service
