# etc/crontab
* *     * * *   root    /home/pi/balkonkraftwerk/check_tibber.py || systemctl restart tibber_pulse.service
